# ocs.service
[Unit]
Description=ocs5 server
Requires=mysql

[Service]
ExecStart=STARTSCRIPT
Type=simple

[Install]
WantedBy=multi-user.target
